# deploy/systemd/control-deck-web.service.in
[Unit]
Description=Ubuntu Control Deck Web
After=network-online.target

[Service]
# アプリが READY=1 を通知し、WATCHDOG=1 を送り続ける（ハング時は自動再起動）
Type=notify
NotifyAccess=main
WatchdogSec=30
WorkingDirectory=@REPO_ROOT@/backend
# PostgreSQL切替時だけ存在する0600の専用設定。値はログやunit本文へ展開しない。
EnvironmentFile=-@REPO_ROOT@/config/database.env
# 起動前にowner／mode／symlink／URL方言を検査する。未作成（SQLite既定）は成功。
ExecStartPre=@REPO_ROOT@/.venv/bin/python -m app.database.cli validate-env-file @REPO_ROOT@/config/database.env
# 待受ホスト/ポートは config.yaml（server.host/port）から読む
ExecStart=@REPO_ROOT@/.venv/bin/python -m app.server
Restart=on-failure
RestartSec=3
# 注: NoNewPrivileges は設定しない。Web ターミナルからの sudo（reboot/apt/systemctl 等）を
# 許可するため（Control Deck はシステム管理ツールであり、これは正当な用途）。
# 保護は認証 + RBAC(terminal.use) + 監査ログで担保する。

[Install]
WantedBy=default.target
