# deploy/systemd/control-deck-web.service.in
[Unit]
Description=Ubuntu Control Deck Web
After=network-online.target

[Service]
# アプリが READY=1 を通知し、WATCHDOG=1 を送り続ける（ハング時は自動再起動）
Type=notify
NotifyAccess=main
WatchdogSec=30
WorkingDirectory=@REPO_ROOT@/backend
# 待受ホスト/ポートは config.yaml（server.host/port）から読む
ExecStart=@REPO_ROOT@/.venv/bin/python -m app.server
Restart=on-failure
RestartSec=3
# セキュリティ強化
NoNewPrivileges=true

[Install]
WantedBy=default.target
